# packages/nrd/nrd-0.0.11.tar.gz/nrd-0.0.11/service/linux/nrd.service
[Unit]
Description=NRD - Vite Dev Server Background Service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i
ExecStart=/usr/bin/python3 -m nrd.nrd
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
