# packages/monitorat/monitorat-0.6.6-py3-none-any.whl/monitorat/systemd/monitor@pip.service
[Unit]
Description=%i Monitor
After=network.target

[Service]
User=__user__
Group=__group__
Environment="PATH=/home/__user__/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/env gunicorn -w 1 -b 0.0.0.0:6161 --timeout 120 monitorat.monitor:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
