# packages/wooting-mouse/wooting_mouse-1.0.2-py3-none-any.whl/data/wooting-mouse@.service
[Unit]
Description=Wooting Mouse daemon

[Service]
Type=simple
ExecStart=/usr/bin/wooting-mouse /%I
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"
