# packages/AWSTransactionalApps/awstransactionalapps-0.15.1.tar.gz/awstransactionalapps-0.15.1/AWSTransactionalAppsStack/scripts/systemd/awsec2-app.service

[Unit]
Description=Run setup.sh at startup
After=network.target

[Service]
Type=oneshot
ExecStart={{setup_script_path}}
RemainAfterExit=true

[Install]
WantedBy=multi-user.target