# packages/surf-controller/surf_controller-1.1.3.tar.gz/surf_controller-1.1.3/scheduler/crontab
# scheduler/crontab
# Run the script every night at 21:00 (9 PM)
# The script's output is redirected to logs.log inside the volume
0 21 * * * python /app/pause_job.py >> /var/log/cron.log 2>&1

# An empty line is required at the end of the file
