# packages/sing-box-config/sing_box_config-0.6.0.tar.gz/sing_box_config-0.6.0/roles/sing_box_install/templates/etc/systemd/system/sing-box-reload.service.j2
[Unit]
Description=Reload sing-box after config changes

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sing-box -C {{ sing_box_etc_dir }} check
ExecStart=/usr/bin/systemctl reload sing-box.service
