# packages/genx-ai/genx_ai-2.10.2-cp312-cp312-manylinux_2_17_x86_64.manylinux2014_x86_64.whl/genx_ai/service/genx-ai.service
[Unit]
Description=GenX AI Service
After=network.target

[Service]
ExecStart=/opt/venv/bin/genx listen
EnvironmentFile=/etc/genx-ai-default.conf
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
